# monitoring/ansible/roles/monitoring/templates/monitoring-stack.service.j2
[Unit]
Description=Monitoring Stack (Prometheus + Grafana + Loki + Promtail)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/monitoring
ExecStartPre=/usr/bin/docker compose -f /opt/monitoring/docker-compose.yml pull --quiet
ExecStart=/usr/bin/docker compose -f /opt/monitoring/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/monitoring/docker-compose.yml down
ExecReload=/usr/bin/docker compose -f /opt/monitoring/docker-compose.yml up -d --remove-orphans
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
